# src/subtree/openbmc/libmctp/systemd/system/mctp-demux.socket
[Unit]
Description=MCTP demux Unix domain socket

[Socket]
ListenSequentialPacket=@mctp-mux
Accept=no

[Install]
WantedBy=sockets.target
